# deploy/jobs.cron
0 8 * * * systemctl restart tradepy.autotrader.service
